# upload/system_tests/conf/test-supervisor.conf
[program:upload-system-tests]
command = python main.py
directory = /app
user = www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
